# etc/profile.d/pk.sh
# pk package manager environment
export PK_ROOT="/"
